# build/SOURCES/supervisor.conf
[program:restify]
command = /appl/restify/venv/bin/gunicorn -c /appl/restify/etc/gunicorn_config.py restify.wsgi
directory = /appl/restify
user = restify
